# capms/context/etc/systemd/system/kubeadm.service.d/kubeadm-load-images.conf
[Service]
ExecStartPre=/usr/local/bin/kubeadm-load-images.sh
