# capms/context/etc/systemd/system/kubeadm.service.d/kubevip-prepare.conf
[Service]
ExecStartPre=/usr/local/bin/kubeadm-prepare-kubevip.sh
