# debian/context/etc/systemd/system/systemd-timesyncd.service.d/override.conf
# Ubuntu sets ConditionVirtualization=!container which prevents integration tests to pass
[Unit]
ConditionVirtualization=