# firewall/context/etc/systemd/system/frr.service.d/override.conf
# Since frr-10 mgmtd requires this directory to be writeable
[Service]
ExecStartPre=/usr/bin/chown -R frr:frr /var/lib/frr
