# cluster-api/context/etc/systemd/system/kubeadm.service.d/boot-order.conf
# kubeadm requires a running containerd.
[Unit]
After=containerd.service
Wants=containerd.service
