# firewall/context/etc/systemd/system/chrony@.service.d/proc.conf
# Without these settings chrony does not start on a firewall.
[Service]
ProcSubset=all
ProtectControlGroups=no