# firewall/context/etc/systemd/system/suricata-update.timer
[Unit]
Description=Run suricata-update hourly and at boot

[Timer]
OnBootSec=10min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target