# debian/context/etc/systemd/system/containerd.service.d/40-ulimit.conf
[Service]
LimitNPROC=infinity
LimitNOFILE=1048576
